# euler-swap-api.service
[Unit]
Description=Euler Swap API
After=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/euler-swap-api

# Environment setup
Environment=NODE_ENV=production
Environment=HOME=/home/admin
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/admin/.local/share/pnpm
Environment=PNPM_HOME=/home/admin/.local/share/pnpm
EnvironmentFile=/home/admin/euler-swap-api/.env

# Start command - using absolute paths
ExecStart=/home/admin/.local/share/pnpm/pnpm start

# Restart policy
Restart=always
RestartSec=10

StandardOutput=append:/var/log/euler-swap-api.log
StandardError=append:/var/log/euler-swap-api.error.log

[Install]
WantedBy=multi-user.target